# deploy/supervisor/ipure-queue.conf
; ============================================================
; IPure Herbs — Supervisor Queue Worker Config
; File: /etc/supervisor/conf.d/ipure-queue.conf
; After placing: supervisorctl reread && supervisorctl update
; ============================================================

[program:ipure-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/backend/artisan queue:work database --tries=3 --timeout=60 --sleep=3 --max-jobs=500 --max-time=3600
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/backend/storage/logs/queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
user=www-data
stopwaitsecs=90
